# launchd/ai.openclaw.qqbot-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.qqbot-bridge</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Volumes/File/proj/qqbot/scripts/start-bridge.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/File/proj/qqbot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/alanz/Library/Logs/qqbot/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alanz/Library/Logs/qqbot/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
  </dict>
</plist>
